# recipes-example/example/datapond-0.1/datapond.service
[Unit]
Description=Datapond daemon

[Service]
Type=simple
ExecStart=/usr/bin/datapond

[Install]
WantedBy=multi-user.target

