# recipes-example/example/collectfw-0.1/collectfw.service
[Unit]
Description=Data collector daemon

[Service]
Type=simple
ExecStart=/usr/bin/collector

[Install]
WantedBy=multi-user.target

